# podman_startup/quadlet/app-bootstrap.container.in
[Unit]
Description=Agentic Chatbot v3 application bootstrap
After=rag-postgres.service ollama-bootstrap.service
Requires=rag-postgres.service ollama-bootstrap.service

[Container]
Image=localhost/agentic-chatbot-v3-app:latest
Pull=never
ContainerName=app-bootstrap
Network=agentic.network
EnvironmentFile=__ENV_DIR__/app.env
Environment=WAIT_FOR_DB=true
Environment=APP_AUTO_MIGRATE=true
Environment=APP_AUTO_INIT_KB=false
Environment=APP_AUTO_INDEX_SKILLS=true
Environment=APP_AUTO_SYNC_KB=true
Environment=DB_WAIT_TIMEOUT_SECONDS=90
Volume=__REPO_ROOT__/data/agents:/app/data/agents:ro,z
Volume=__REPO_ROOT__/data/prompts:/app/data/prompts:ro,z
Volume=__REPO_ROOT__/data/skills:/app/data/skills:ro,z
Volume=__REPO_ROOT__/data/skill_packs:/app/data/skill_packs:ro,z
Volume=__REPO_ROOT__/data/router:/app/data/router:ro,z
Volume=__REPO_ROOT__/data/kb:/app/data/kb:ro,z
Volume=__REPO_ROOT__/data/demo:/app/data/demo:ro,z
Volume=__REPO_ROOT__/docs:/app/docs:ro,z
Volume=app_runtime.volume:/app/data/runtime
Volume=app_uploads.volume:/app/data/uploads
Volume=app_workspaces.volume:/app/data/workspaces
Volume=app_cache.volume:/app/data/cache
Volume=app_control_panel_data.volume:/app/data/control_panel
Volume=graphrag_projects.volume:/app/data/graphrag/projects
Exec=sh -c "python run.py runtime-smoke --registry-only --json && echo 'Application bootstrap complete'"

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=0

[Install]
WantedBy=default.target
